# roles/thunderbird_headless/templates/thunderbird-headless.service.j2
[Unit]
Description=Thunderbird (headless) - background mail & calendar sync
After=network-online.target dbus.service
Wants=network-online.target

[Service]
Type=simple
Environment=MOZ_HEADLESS=1
Environment=HOME=%h
Environment=XDG_RUNTIME_DIR=%t
WorkingDirectory=%h
{% if thunderbird_use_xvfb -%}
ExecStart={{ thunderbird_xvfb_run }} {{ thunderbird_exec }} --headless{% if thunderbird_extra_args|length %} {{ thunderbird_extra_args }}{% endif %}
{% else -%}
ExecStart={{ thunderbird_exec }} --headless{% if thunderbird_extra_args|length %} {{ thunderbird_extra_args }}{% endif %}
{% endif -%}
# Restart policy:
{% if thunderbird_refresh_method == 'runtime_max' -%}
Restart=always
{% else -%}
Restart=on-failure
{% endif -%}
RestartSec=5s
TimeoutStopSec=10s
{% if thunderbird_refresh_method == 'runtime_max' -%}
# Auto-restart the service after this runtime (forces periodic refresh)
RuntimeMaxSec={{ thunderbird_refresh_interval }}
{% endif -%}

[Install]
WantedBy=default.target

