# roles/thunderbird_headless/templates/thunderbird-headless-refresh.service.j2
[Unit]
Description=Restart Thunderbird headless to force sync

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user restart {{ thunderbird_service_name }}

