# roles/mcp_hub_ui/templates/mcp-hub-ui.service.j2
# roles/mcp_hub_ui/templates/mcp-hub-ui.service.j2
[Unit]
Description=MCP Hub UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart={{ mcp_hub_ui_bin }} ui --port {{ mcp_hub_ui_port }}{% for a in mcp_hub_ui_extra_args %} {{ a }}{% endfor %}
Restart=on-failure
RestartSec=2
{% for k, v in mcp_hub_ui_environment.items() %}
Environment={{ k }}={{ v }}
{% endfor %}
WorkingDirectory=%h

[Install]
WantedBy=default.target

